# launchd/com.qi.claude-cache-status.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.qi.claude-cache-status</string>

    <key>ProgramArguments</key>
    <array>
      <string>__STATUS_APP_PATH__</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__/menubar-app</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>CLAUDE_CACHE_STATUS_URL</key>
      <string>__CLAUDE_CACHE_STATUS_URL__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/claude-cache-status.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/claude-cache-status.error.log</string>
  </dict>
</plist>
